# containers/nginx/letsencrypt_retrieve_keys.cron
SHELL=/bin/bash
0 0 * * 1 /usr/bin/env bash /letsencrypt_retrieve_keys.sh
